# systemd/InoSwitchboardController.service
[Unit]
Description=Arduino Switchboard Control Service

[Service]
CPUQuota=5%

Type=simple
ExecStart=%h/.local/bin/SwitchboardController
